# .platform/cron/renew_cert.cron
0 0,12 * * * root certbot renew --quiet --post-hook "systemctl restart nginx"
